# scraper/com.csekeyan.fifa-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.csekeyan.fifa-scraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/kgnmzn/Desktop/AI/Experiments/FIFA/scraper/scrape.py</string>
        <string>--once</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>WorkingDirectory</key>
    <string>/Users/kgnmzn/Desktop/AI/Experiments/FIFA</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WORKER_URL</key>
        <string>https://fifa-wc-2026.pages.dev/api/results</string>
        <key>API_KEY</key>
        <string>fifa-wc-2026-scraper-key</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kgnmzn/Desktop/AI/Experiments/FIFA/scraper/scraper.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kgnmzn/Desktop/AI/Experiments/FIFA/scraper/scraper.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
